# rrrr.conf
[program:rrrr]
command=python3 /home/curator2/rrrr/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/rrrr/test.err.log
stdout_logfile=/home/curator2/rrrr/test.out.log